# systemd/school-mosreg.timer
[Unit]
Description=Checking the school schedule timer

[Timer]
OnCalendar=*-*-* 16:05:00
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
